# tomcat.conf
[supervisord]
nodaemon=true

[program:tomcat8]
command=/usr/sbin/tomcat8 start
user=tomcat
